# scripts/install/support-files/templates/#etc#supervisor-consul.conf
[unix_http_server]
file=__BK_HOME__/logs/consul/supervisor.sock   ; (the path to the socket file)

;[inet_http_server]         ; inet (TCP) server disabled by default
;port=127.0.0.1:9001        ; (ip_address:port specifier, *:port for all iface)
;username=user              ; (default is no username (open server))
;password=123               ; (default is no password (open server))

[supervisord]
logfile=__BK_HOME__/logs/consul/supervisord_out.log ; (main log file;default $CWD/supervisord_out.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=__BK_HOME__/logs/consul/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://__BK_HOME__/logs/consul/supervisor.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as http_username if set
;password=123                ; should be same as http_password if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available

[program:consul-agent]
directory=/usr/bin/
command=/usr/bin/consul agent -config-file=__BK_HOME__/etc/consul.conf -config-dir=__BK_HOME__/etc/consul.d
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=__BK_HOME__/logs/consul/consul.log
